# install/ssh-summarize.service
[Unit]
Description=Summarize new SSH honeypot sessions

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=/root/ximg-web/ssh-server/summarize-sessions.sh
WorkingDirectory=/root/ximg-web
StandardOutput=append:/root/ximg-web/ssh-logs/summarize.log
StandardError=append:/root/ximg-web/ssh-logs/summarize.log
