# install/ssh-summarize.path
[Unit]
Description=Watch ssh-logs for new honeypot session files

[Path]
PathModified=/root/ximg-web/ssh-logs
Unit=ssh-summarize.service

[Install]
WantedBy=multi-user.target
